# custom-nginx/rootfs/etc/s6-overlay/s6-rc.d/crond/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start crond service
# Source: https://github.com/home-assistant/addons/tree/master/nginx_proxy
# ==============================================================================

exec crond -f -d 9
